# roles/healthcheck_heartbeat/templates/healthcheck-heartbeat.timer.j2
[Unit]
Description=Send regular updates to healthchecks.io

[Timer]
OnCalendar={{ healthcheck_heartbeat_time_period }}
Persistent=true

[Install]
WantedBy=timers.target
